# tests/auto/network/support/modules/network_test_server/files/cron.daily/writeableswatch
#!/bin/sh
exec /usr/sbin/tmpreaper 1d /home/writeables/dav/ /home/writeables/smb/ /home/qt-test-server/ftp/qtest/upload/
